# ops/systemd/idox-worker.service
[Unit]
Description=Idox scrape_jobs worker (LISTEN/claim/execute scraper)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/scraper
Environment=NODE_ENV=production
EnvironmentFile=-/etc/default/idox-worker
# Recommended: run as a non-root user that has Playwright browsers installed (or can install them).
#User=james
#Group=james
ExecStart=/usr/bin/node /opt/scraper/worker_listen_new.js
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
