# ops/systemd/idox-discovery.timer
[Unit]
Description=Run Idox validated-range discovery hourly (oneshot)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

